# infra/systemd/iseelocal-relay.service
[Unit]
Description=iseelocal relay
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=iseelocal
Group=iseelocal
EnvironmentFile=/etc/iseelocal/relay.env
ExecStart=/usr/local/bin/iseelocal-relay
Restart=always
RestartSec=3
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ProtectHome=true
ReadWritePaths=/var/lib/iseelocal

[Install]
WantedBy=multi-user.target
